# config/systemd/o4o-webserver.service
[Unit]
Description=O4O Platform Web Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/home/ubuntu/o4o-platform
Environment="NODE_ENV=production"
Environment="PATH=/usr/bin:/usr/local/bin"

# Main Site (포트 3000)
ExecStart=/usr/bin/serve -s /home/ubuntu/o4o-platform/apps/main-site/dist -l 3000

# 재시작 정책
Restart=always
RestartSec=10
StandardOutput=append:/var/log/o4o-webserver.log
StandardError=append:/var/log/o4o-webserver.error.log

[Install]
WantedBy=multi-user.target